# meta-gap/recipes-apps/ntp/files/ntpdate.crontab
*/5 * * * * root /usr/bin/ntpdate-sync silent
